# Mini-Prj-1/service/result.service
[Unit]
Description=Service for voting-app result

[Service]
ExecStart=/opt/devops/Mini-Prj-1/runner/result.sh
Restart=always

[Install]
WantedBy=multi-user.target