# Mini-Prj-1/service/vars.service
[Unit]
Description=Service for voting-app initializing

[Service]
ExecStart=/opt/devops/Mini-Prj-1/runner/proj-vars.sh
Type=oneshot
# Restart=always

[Install]
WantedBy=multi-user.target